# amdgpu-zerodb.service
[Unit]
Description=0dB service for AMDGPU

[Service]
Type=simple
ExecStart=/usr/bin/amdgpu-zerodb

[Install]
WantedBy=multi-user.target
